# com.classic-hunt.scrape.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.classic-hunt.scrape</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/mariokarlovcec/Learn/github/classic-hunt/run-scrape.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mariokarlovcec/Learn/github/classic-hunt</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/mariokarlovcec/Learn/github/classic-hunt/logs/scrape.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mariokarlovcec/Learn/github/classic-hunt/logs/scrape-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
